# services/ingestion-service/deployments/batch/crontab
# Ingestion Service Batch Processing Schedule
# Format: minute hour day month weekday command

# Collect trending videos - Twice daily at 3:00 AM and 3:00 PM
0 3,15 * * * /app/bin/batch-trending

# Schedule snapshots for recent videos - Every hour
0 * * * * /app/bin/batch-schedule-snapshots -hours 2

# Generate rankings - Daily at 6:00 AM
0 6 * * * /app/bin/batch-rankings

# Renew WebSub subscriptions - Daily at 1:00 AM
0 1 * * * /app/bin/batch-websub-renewal -days 7

# Note: All times are in UTC
# Adjust according to your target timezone and peak usage patterns